# newsboat/.config/systemd/user/newsboat-oneshot.timer
[Unit]
Description=Newsboat Timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=1h

[Install]
WantedBy=default.target
